# profile.d/arcanist.sh
# arcanist
ARC_PATH=~/dev/oss/arcanist/bin
if [ -d ${ARC_PATH} ]; then
    export PATH=${ARC_PATH}:${PATH}
fi

