# profile.d/heroku.sh
# heroku toolbelt
HEROKU_PATH=/usr/local/heroku
if [ -d ${HEROKU_PATH} ]; then
    export PATH=${HEROKU_PATH}/bin:$PATH
fi

