# profile.d/node.sh
#NodeJS
NODEROOT=~/dev/nodejs
if [ -d ${NODEROOT} ]; then
   export PATH=${NODEROOT}/bin:$PATH
fi
