# profile.d/ragel.sh
# ragel
RAGEL_PATH=~/dev/ragel-6.9
if [ -d ${RAGEL_PATH} ]; then
    export PATH=${RAGEL_PATH}/bin:$PATH
fi
